# service/pgsqlsnapshot/run
#!/bin/sh

export HOME=/var/lib/postgresql
. $HOME/v/.prod/local/bin/activate

export PATH=$PWD:$PATH

exec 2>&1
exec setuidgid postgres account_position.py
